# systemd-units/node-example.socket
[Socket]
ListenStream=3333

[Install]
WantedBy=sockets.target
